# gateway/echosmart-gateway.service
[Unit]
Description=EchoSmart Gateway Service
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
SupplementaryGroups=gpio i2c spi dialout
EnvironmentFile=/etc/echosmart/gateway.env
ExecStart=/usr/bin/echosmart-gateway run
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=echosmart-gateway

[Install]
WantedBy=multi-user.target
