# gateway/debian/echosmart-gateway.service
[Unit]
Description=EchoSmart Gateway Service
Documentation=https://github.com/Uaemextop/echosmart
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=echosmart
Group=echosmart
WorkingDirectory=/opt/echosmart/gateway
ExecStart=/usr/bin/echosmart-gateway run
Restart=on-failure
RestartSec=10
WatchdogSec=120

# Environment
EnvironmentFile=-/etc/echosmart/gateway.env
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=echosmart-gateway

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/echosmart /opt/echosmart/gateway
PrivateTmp=true

[Install]
WantedBy=multi-user.target
