# gateway/debian/prerm
#!/bin/bash
# prerm — runs before echosmart-gateway is removed

set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop and disable the service
        if systemctl is-active --quiet echosmart-gateway.service 2>/dev/null; then
            systemctl stop echosmart-gateway.service
        fi
        systemctl disable echosmart-gateway.service 2>/dev/null || true
        ;;
esac

#DEBHELPER#

exit 0
